# DDBOT-WSa-Desktop/service/ddbot-webui.service
[Unit]
Description=DDBOT WebUI Service
After=network.target

[Service]
Type=simple
User=ddbot
WorkingDirectory=/opt/ddbot-webui/backend
ExecStart=/opt/ddbot-webui/backend/target/release/ddbot-backend
Restart=always
RestartSec=10

# 环境变量
Environment=RUST_LOG=info
Environment=PORT=3000

# 日志配置
StandardOutput=journal
StandardError=journal

# 安全设置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target